# Configs/tools-dashboard.ini
[program:td-discovery]
command=/usr/local/jdk/bin/java -jar -Xms64m -Xmx312m tools-dashboard-discovery.jar
directory=/home/webapp/
priority=1
user=root
autostart=true
startretries=3
autorestart=true

[program:td-config]
command=/usr/local/jdk/bin/java -jar -Xms64m -Xmx312m tools-dashboard-config.jar
directory=/home/webapp/
priority=1
user=root
autostart=true
startretries=3
autorestart=true

[program:td-api-gateway]
command=/usr/local/jdk/bin/java -jar -Xms64m -Xmx312m tools-dashboard-api-gateway.jar
directory=/home/webapp/
priority=1
user=root
autostart=true
startretries=3
autorestart=true

[program:td-authorization]
command=/usr/local/jdk/bin/java -jar -Djava.util.concurrent.ForkJoinPool.common.parallelism=50 tools-dashboard-authorization.jar
directory=/home/webapp/
user=root
autostart=true
autorestart=true

[program:td-baidu]
command=/usr/local/jdk/bin/java -jar -Xms64m -Xmx512m -Djava.util.concurrent.ForkJoinPool.common.parallelism=26 tools-dashboard-baidu.jar
directory=/home/webapp/
user=root
autostart=true
autorestart=true

[program:td-aatool]
command=/usr/local/jdk/bin/java -jar -Djava.util.concurrent.ForkJoinPool.common.parallelism=50 tools-dashboard-aatool.jar
directory=/home/webapp/
user=root
autostart=true
autorestart=true

[program:td-aplus]
command=/usr/local/jdk/bin/java -jar -Dco.paralleluniverse.fibers.detectRunawayFibers=false tools-dashboard-myspace-cap.jar
directory=/home/webapp/
user=root
autostart=true
startretries=1
autorestart=true

[program:td-web-facade]
command=/usr/local/jdk/bin/java -jar -Djava.util.concurrent.ForkJoinPool.common.parallelism=16 tools-dashboard-web-facade.jar
directory=/home/webapp/
user=root
autostart=true
startretries=1
autorestart=true
